# cookbooks/elastic-search/templates/default/elastic-search.conf.upstart
description 'Elastic Search Server'

start on startup
stop on shutdown

setuid '__UID__'
setgid '__GID__'

respawn
respawn limit 10 5

env JAVA_HOME='__JDK_FOLDER__'

exec '__INSTALL_FOLDER__/bin/elasticsearch' -Des.config='__INSTALL_FOLDER__/config/elasticsearch.yml'

post-stop script
    rm --force --recursive __INSTALL_FOLDER__/logs/*
end script