# cookbooks/docker/files/default/aufs.conf.upstart
description 'AUFS Module'

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

limit nofile 524288 1048576
limit nproc 524288 1048576

kill timeout 20

script
    lsmod | awk '{ print $1 }' | grep -F -o 'aufs' > '/dev/null' || modprobe -q 'aufs' || apt-get install -y "linux-image-extra-$(uname -r)" && modprobe 'aufs'
end script