# cookbooks/mongodb/files/upstart/mongodb.conf
description 'MongoDB Server'

start on startup
stop on shutdown

exec '/opt/mongodb/bin/mongod' --port 27017 --dbpath '/opt/mongodb/data'

respawn
respawn limit 10 5
