# cookbooks/tomcat/templates/default/tomcat.conf.upstart
description 'Tomcat Server'

start on started mountall
stop on shutdown

setuid '__USER_NAME__'
setgid '__GROUP_NAME__'

respawn
respawn limit 10 5

env CATALINA_HOME='__INSTALL_FOLDER__'
env HOME='__HOME_FOLDER__'
env JAVA_HOME='__JDK_INSTALL_FOLDER__'

exec "${CATALINA_HOME}/bin/catalina.sh" run

post-stop script
    rm -f -r ${CATALINA_HOME}/logs/* ${CATALINA_HOME}/temp/*
end script