# cookbooks/mongodb/templates/default/mongodb.conf.conf
description 'MongoDB Server'

start on startup
stop on shutdown

respawn
respawn limit 10 5

exec '__INSTALL_FOLDER__/bin/mongod' --port __PORT__ --dbpath '__INSTALL_DATA_FOLDER__'