# cookbooks/mongodb/templates/default/mongodb.conf.upstart
description 'MongoDB Server'

start on started mountall
stop on shutdown

setuid '__USER_NAME__'
setgid '__GROUP_NAME__'

respawn
respawn limit 10 5

exec '__INSTALL_FOLDER__/bin/mongod' --port __PORT__ --dbpath '__INSTALL_DATA_FOLDER__'