# cookbooks/tomcat/templates/tomcat.cron
#!/bin/bash -e

chown -R '__USER_NAME__:__GROUP_NAME__' '__INSTALL_FOLDER_PATH__'