# cookbooks/go-server/files/upstart/go-server.conf
description 'Go-Server (Server)'

start on startup
stop on shutdown

setuid '__UID__'
setgid '__GID__'

respawn
respawn limit 10 5

chdir '__SERVER_INSTALL_FOLDER__'
exec java -jar go.jar