# cookbooks/haproxy/templates/default/haproxy.conf.upstart
description 'HAProxy Server'

start on startup
stop on shutdown

env DAEMON='__INSTALL_FOLDER__/sbin/haproxy'

expect fork

respawn
respawn limit 10 5

exec $DAEMON

pre-start script
    $DAEMON -t

    if [ $? -ne 0 ]
        then exit $?
    fi
end script

post-stop script
    rm -f -r __INSTALL_FOLDER__/logs/* __INSTALL_FOLDER__/*_temp
end script