# cookbooks/haproxy/templates/default/haproxy.conf.upstart
description 'HAProxy Server'

start on started mountall
stop on shutdown

expect fork

respawn
respawn limit 10 5

exec '__INSTALL_FOLDER__/sbin/haproxy'