# cookbooks/ghost/templates/default/ghost.conf.upstart
description 'Ghost Server'

start on startup
stop on shutdown

setuid '__UID__'
setgid '__GID__'

respawn
respawn limit 10 5

chdir '__INSTALL_FOLDER__'
exec npm start