# cookbooks/ghost/templates/default/ghost.conf.upstart
description 'Ghost Server'

start on startup
stop on shutdown

setuid '__UID__'
setgid '__GID__'

respawn
respawn limit 10 5

env NODE_ENV='__ENVIRONMENT__'

exec node '__INSTALL_FOLDER__/index'