# cookbooks/selenium-server/templates/default/selenium-server-node.conf.upstart
description 'Selenium Server (Node)'

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown

setuid '__USER_NAME__'
setgid '__GROUP_NAME__'

respawn
respawn limit 10 5

exec java -jar '__INSTALL_FILE__' -role 'node' -trustAllSSLCertificates -nodeConfig '__CONFIG_FILE__'