# cookbooks/selenium-server/templates/default/selenium-server.conf.upstart
description 'Selenium Server'

start on started mountall
stop on shutdown

setuid '__USER_NAME__'
setgid '__GROUP_NAME__'

respawn
respawn limit 10 5

exec java -jar '__INSTALL_FILE__' -role 'hub'