# cookbooks/redis/files/upstart/redis.conf
description 'Redis Server'

start on startup
stop on shutdown

setuid 'redis'
setgid 'redis'

exec '/opt/redis/bin/redis-server' '/opt/redis/config/redis.conf'

respawn
respawn limit 10 5
