# cookbooks/ghost/files/upstart/ghost.conf
description 'Ghost Server'

start on startup
stop on shutdown

setuid '__UID__'
setgid '__GID__'

respawn
respawn limit 10 5

exec node '__INSTALL_FOLDER__/index'
