# cookbooks/ghost/files/upstart/ghost.conf
description 'Ghost Server'

start on startup
stop on shutdown

exec node '/opt/ghost/index'

respawn
respawn limit 10 5
