# cookbooks/nginx/templates/default/nginx.conf.upstart
description 'Nginx Server'

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown

limit nofile 524288 1048576
limit nproc 524288 1048576

kill timeout 20

respawn
respawn limit 10 5

expect fork

exec '__INSTALL_FOLDER__/sbin/nginx'

post-stop script
    rm -f -r __INSTALL_FOLDER__/*_temp
    cd "__INSTALL_FOLDER__/logs" && find '.' -not -name '.' -delete
end script