# cookbooks/nginx/templates/default/nginx.conf.upstart
description 'Nginx Server'

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown

expect fork

respawn
respawn limit 10 5

exec '__INSTALL_FOLDER__/sbin/nginx'

post-stop script
    rm -f -r __INSTALL_FOLDER__/*_temp
    cd "__INSTALL_FOLDER__/logs" && find '.' -not -name '.' -delete
end script