# cookbooks/nginx/templates/default/nginx.conf.upstart
description 'Nginx Server'

start on started mountall
stop on shutdown

expect fork

respawn
respawn limit 10 5

exec '__INSTALL_FOLDER__/sbin/nginx'

post-stop script
    rm -f -r __INSTALL_FOLDER__/*_temp
    cd "__INSTALL_FOLDER__/logs" && find '.' ! -name '.' -delete
end script